# os customizing/etc/cron.d/geoip
# Jeden 8. des Monats update der GEOIP Datenbank
0 0 8 1 * root /opt/bi-s/software/scripts/elk/geoip-update.sh >> /var/log/geoip-update.log 2>&1