# os customizing/etc/cron.d/osiris
#License check um 03:00
00 03 * * * root /opt/bi-s/software/scripts/lic/check_license.sh 2>&1
30 03 * * * root /opt/bi-s/software/scripts/lic/shutdown 2>&1