# os customizing/etc/cron.d/osiris
#License check jede Stunde
0 */1 * * * root /opt/bi-s/software/scripts/lic/check_license 2>&1
30 03 * * * root /opt/bi-s/software/scripts/lic/shutdown 2>&1