# os customizing/etc/cron.d/osiris

#License check jeden Tag um 01:00
0 1 * * * root /opt/bi-s/software/scripts/lic2/check_lic 2>&1