# os customizing/etc/cron.d/sync_poller
#Poller wird jede Stunde syncronisiert
0 */1 * * * root /opt/bi-s/software/scripts/sync/sync_poller.sh