# os customizing/etc/profile.d/oracle.sh
ORACLE_BASE=/opt/oracle
ORACLE_HOME=$ORACLE_BASE/112
LD_LIBRARY_PATH=$ORACLE_HOME/lib
PATH=$PATH:$ORACLE_HOME/bin

export ORACLE_BASE ORACLE_HOME LD_LIBRARY_PATH PATH